# recipes-core/lsbscripts/files/brup.service
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/zipgw-scripts/ifup.d/50bridge.sh

[Install]
WantedBy=multi-user.target


