# recipes-core/lsbscripts/files/zipgateway.service
[Unit]
Description=LSB: start Z/IP Gateway (zipgateway)
After=systemd-networkd
Conflicts=shutdown.target

[Service]
Type=forking
ExecStart=/etc/init.d/zipgateway start
ExecStop=/etc/init.d/zipgateway stop
RestartSec=3
Restart=on-failure

[Install]
WantedBy=multi-user.target

